# stock_controller/crontab
* * * * * python /app/src/collecter/Stock.py request -symbol MSFT


